# src/static/static/home/yi/etc/crontabs/root
12 2 * * * /home/yi/script/check_update.sh
